# docker-image/config/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apachectl -DFOREGROUND
[program:redis]
command=redis-server
autorestart=true

[program:mysql]
command=/bin/sh /tmp/run.sh
